# ops/deepapi-backup.timer
[Unit]
Description=Run DeepAPI encrypted offsite backup daily

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=15m
Unit=deepapi-backup.service

[Install]
WantedBy=timers.target
